# deploy/ansible/roles/snapshotter/files/goshimmer-snapshot.bin-dump.service
[Unit]
Description=Dump ZIPP Snapshots

[Service]
Type=oneshot
ExecStart=/opt/zipp_snapshot.bin_dump.sh

[Install]
WantedBy=timers.target
